# systemd/llm-wiki-mcp.service
[Unit]
Description=llm_wiki MCP stdio server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# 先复制 config/llm_wiki.template.toml 为本地 config/llm_wiki.toml，并修改 knowledge_root
WorkingDirectory=/path/to/llm_wiki
Environment=ORT_DYLIB_PATH=/path/to/onnxruntime/capi/libonnxruntime.so
Environment=LD_LIBRARY_PATH=/path/to/onnxruntime/capi
ExecStart=/path/to/llm_wiki/target/release/llm-wiki --config /path/to/llm_wiki/config/llm_wiki.toml serve-mcp
Restart=on-failure

[Install]
WantedBy=multi-user.target
